# BuildSupport/dev/com.shin.cellcap.helper.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shin.cellcap.helper</string>

    <key>MachServices</key>
    <dict>
        <key>com.shin.cellcap.helper</key>
        <true/>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>__HELPER_PATH__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__STDOUT_PATH__</string>

    <key>StandardErrorPath</key>
    <string>__STDERR_PATH__</string>
</dict>
</plist>
